# application-followups.service
[Unit]
Description=Application Followups Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/application-followups
ExecStart=/root/.local/share/pnpm/pnpm start
Restart=always
Environment=NODE_ENV=production
Environment=PORT=8000

[Install]
WantedBy=multi-user.target
